# deploy/systemd/webcam-watchdog.timer.example
[Unit]
Description=Run webcam watchdog every minute

[Timer]
OnBootSec=30s
OnUnitActiveSec=60s
Unit=webcam-watchdog.service

[Install]
WantedBy=timers.target
